# DEBIAN/postinst
#!/bin/sh
systemctl enable draginofwd
systemctl start draginofwd
